# scripts/deploy/spinalfmriprep-dashboard.service
[Unit]
Description=SpinalfMRIprep QC Dashboard Server (latest wf)
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
Type=simple
EnvironmentFile=/etc/spinalfmriprep-dashboard.env
ExecStart=/usr/bin/python3 /opt/spinalfmriprep/spinalfmriprep_dashboard_server.py
Restart=always
RestartSec=2
TimeoutStartSec=10

# Optional hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target




